# quadlet/p-php.container
[Unit]
Description=PHP Container

[Container]
Image=p-php.build
ContainerName=php
ServiceName=pc-php

[Install]
WantedBy=multi-user.target
